# provisioning/roles/main/templates/mount-bindfs.conf
# Remount directories with bindfs
description "Remount vagrant shared folder"

start on vagrant-mounted

script
  bindfs -u www-data -g www-data -p u=rwX:g=rD:o=rD /vagrant /var/www/html
end script
